# scripts/linux/football-cups-db-import.service
[Unit]
Description=Football Cups PostgreSQL import
Wants=network-online.target
After=network-online.target postgresql.service

[Service]
Type=oneshot
User=football-cups
Group=football-cups
WorkingDirectory=/opt/football-cups
EnvironmentFile=-/etc/football-cups/collector.env
ExecStart=/opt/football-cups/.venv/bin/football-cups-db import-files --workspace /opt/football-cups
TimeoutStartSec=180
Restart=on-failure
RestartSec=300
StartLimitIntervalSec=1800
StartLimitBurst=3
RestartPreventExitStatus=2

[Install]
WantedBy=multi-user.target
