# scripts/linux/football-cups-db-import.service
[Unit]
Description=Football Cups PostgreSQL import
Wants=network-online.target postgresql.service
After=network-online.target postgresql.service
RequiresMountsFor=/srv/football-cups
ConditionPathIsMountPoint=/srv/football-cups
StartLimitIntervalSec=1800
StartLimitBurst=4

[Service]
Type=oneshot
User=football-cups
Group=football-cups
WorkingDirectory=/opt/football-cups
EnvironmentFile=-/etc/football-cups/collector.env
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/opt/football-cups/.venv/bin/football-cups-db import-files --workspace /opt/football-cups
TimeoutStartSec=180
Restart=on-failure
RestartSec=300
RestartPreventExitStatus=2
UMask=0027
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadOnlyPaths=/opt/football-cups
ReadWritePaths=/srv/football-cups

[Install]
WantedBy=multi-user.target
